# rc.local
# Add to /etc/rc.local

# start pias script
/usr/bin/python3 /home/pi/pias/pias.py &>> /var/log/pias.log